# src/server/cassh.service
# /lib/systemd/system/cassh.service
[Unit]
Description=cassh
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/cassh/server/server.py -c /etc/cassh/cassh.conf
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target

